# roles/githubixx.etcd/templates/etc/systemd/system/etcd.service.j2
{% macro cluster_hosts() %}
{%- for host in groups['etcd_nodes'] -%}
{{hostvars[host]['ansible_hostname']}}=http://{{ hostvars[host].ansible_host }}:{{etcd_peer_port}}{% if not loop.last %},{% endif %}
{%- endfor -%}
{% endmacro %}

[Unit]
Description=etcd
Documentation=https://github.com/coreos

[Service]
ExecStart={{etcd_bin_dir}}/etcd \
  --name {{ansible_hostname}} \
  --advertise-client-urls http://{{hostvars[inventory_hostname]['ansible_host']}}:{{etcd_client_port}} \
  --listen-client-urls http://0.0.0.0:{{etcd_client_port}} \
  --initial-advertise-peer-urls http://{{hostvars[inventory_hostname]['ansible_host']}}:{{etcd_peer_port}} \
  --listen-peer-urls http://0.0.0.0:{{etcd_peer_port}} \
  --initial-cluster-token {{etcd_initial_cluster_token}} \
  --initial-cluster {{cluster_hosts()}} \
  --initial-cluster-state {{etcd_initial_cluster_state}} \
  --data-dir={{etcd_data_dir}}
Restart=on-failure
RestartSec=5
Type=notify

[Install]
WantedBy=multi-user.target
